# dropbox/dropboxperso
# /etc/init.d/dropboxperso
### BEGIN INIT INFO
# Provides:          dropboxperso
# Required-Start:    $local_fs $remote_fs $network $syslog $named
# Required-Stop:     $local_fs $remote_fs $network $syslog $named
# Should-Start:      $named $time
# Should-Stop:       $named $time
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     false
# Short-Description: Start and stop the dropbox daemon for debian/ubuntu
# Description:       Dropbox daemon for linux
### END INIT INFO

USER="david"
HOMEDIR=`getent passwd $USER | cut -d: -f6`
DROPBOX_HOME=.dropbox_perso
DAEMON="$HOMEDIR/$DROPBOX_HOME/.dropbox-dist/dropbox"
PROCESS=$HOMEDIR/$DROPBOX_HOME/.dropbox-dist/dropbox
export HOME=$HOME/$DROPBOX_HOME

start() {
  echo "Starting dropbox perso..."
  if [ -x $DAEMON ]; then
  	start-stop-daemon -v -b -o -c $USER -S -x $DAEMON
  fi
}

stop() {
  echo "Stopping dropbox perso..."
  if [ -x $DAEMON ]; then
  	start-stop-daemon -o -c $USER -K -x $PROCESS
  fi
}

status() {
  dbpid=`ps aux | grep $PROCESS | grep -v grep | head -n 1 | awk '{print $2}'`
	if [ -z $dbpid ] ; then
		echo "dropboxd perso for USER $USER: not running."
	else
		echo "dropboxd perso for USER $USER: running (pid $dbpid)"

	fi
}

case "$1" in
  start)
    start
    ;;

  stop)
    stop
    ;;

  restart|reload|force-reload)
    stop
    start
    ;;

  status)
    status
    ;;

  *)
    echo "Usage: /etc/init.d/dropboxperso {start|stop|reload|force-reload|restart|status}"
    exit 1

esac

exit 0